# openclaw-self-healing/systemd/openclaw-gateway.service
# OpenClaw Gateway - Level 0 KeepAlive (user-level systemd)
# Install: ~/.config/systemd/user/openclaw-gateway.service
#
# NOTE: OpenClaw may already manage its own systemd service via
# `openclaw gateway start`. This unit is provided as a FALLBACK
# for manual installations where OpenClaw doesn't auto-register.
#
# Check first: systemctl --user list-units | grep openclaw
# If OpenClaw already has a service, you can skip this file.

[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'exec node $(npm root -g)/openclaw/dist/gateway/index.js'
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PATH=%h/.local/bin:%h/.nvm/versions/node/current/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
Environment=OPENCLAW_GATEWAY_PORT=18789
WorkingDirectory=%h/openclaw

[Install]
WantedBy=default.target
